# supervisor/conf/apm.ini
[program:apm]
command=/bin/bash -c "/usr/local/lib/apm/apm-server run -e --path.config=/usr/local/lib/apm --strict.perms=false"
process_name=%(program_name)s
numprocs=1
autorestart=false
autostart=true
