# supervisor/conf/elasticsearch.ini
[program:elasticsearch]
command=/bin/su elasticsearch /bin/bash -c "/usr/local/lib/elasticsearch/bin/elasticsearch -p /opt/elasticsearch/system/process.pid"
process_name=%(program_name)s
pidfile=/opt/elasticsearch/system/process.pid
numprocs=1
autorestart=false
autostart=true