# supervisor/conf/kibana.ini
[program:kibana]
command=/bin/su kibana /bin/bash -c "/usr/local/lib/kibana/bin/kibana"
process_name=%(program_name)s
numprocs=1
autorestart=false
autostart=true
