# supervisor/conf/logstash.ini
[program:logstash]
command=/bin/su logstash /bin/bash -c "/usr/local/lib/logstash/bin/logstash"
process_name=%(program_name)s
numprocs=1
autorestart=false
autostart=true